# archive-bot.service
[Unit]
Description=Discord Archive Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/kasap-bot
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
